# nts-schedule-notifier.service
# Systemd service file for NTS Schedule Notifier
# 
# Installation:
#   1. Copy this file to /etc/systemd/system/
#   2. Edit ExecStart path and User as needed
#   3. Run: sudo systemctl daemon-reload
#   4. Run: sudo systemctl enable nts-schedule-notifier.service
#
# Note: This service is designed to be triggered by a systemd timer,
# not run continuously. See nts-schedule-notifier.timer for the timer config.

[Unit]
Description=NTS Radio Daily Schedule Notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/pi/.local/bin/uv run /home/pi/nts-daily-schedule-notifier/nts_schedule_notifier.py
WorkingDirectory=/home/pi/nts-daily-schedule-notifier
User=pi
Environment="PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin"
# Note: NTFY_TOPIC and NTFY_SERVER are loaded from .env file in the project directory

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
